# deploy/systemd/asl-cactus-hybrid.service
[Unit]
Description=ASL Cactus Hybrid Inference Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=asl
Group=asl
WorkingDirectory=/opt/sign-language-asl
EnvironmentFile=/etc/asl/cactus-hybrid.env
ExecStart=/usr/bin/python3 /opt/sign-language-asl/scripts/runtime/serve_cactus_hybrid_service.py --host 0.0.0.0 --port 9000
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
